# systemd/ozzybot-telegram-cmd.service
[Unit]
Description=OzzyBot Telegram Command Bot (Jarvis Interface)
After=network-online.target ozzybot-webhook.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=/home/rick/ozzy-bot
ExecStart=/home/rick/ozzy-bot/venv/bin/python /home/rick/ozzy-bot/telegram_command_bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=default.target
