# systemd/ozzybot-orphan-report.timer
[Unit]
Description=Run OzzyBot Hermes orphan reconciliation report daily at 06:00
After=network-online.target
Wants=network-online.target

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
