# systemd/ozzybot-testnet-15m-reversal.timer
[Unit]
Description=Run OzzyBot Testnet 15m Reversal Capture Scanner every 15 minutes

[Timer]
OnCalendar=*:02/15
Persistent=true
Unit=ozzybot-testnet-15m-reversal.service

[Install]
WantedBy=timers.target
