# systemd/ozzybot-openclaw-market-sniper.timer
[Unit]
Description=Run OzzyBot OpenClaw Market Sniper every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=ozzybot-openclaw-market-sniper.service

[Install]
WantedBy=timers.target
