# systemd/ozzybot-signal.service
[Unit]
Description=OzzyBot Hermes 1h signal generator
After=network-online.target ozzybot-webhook.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/rick/ozzy-bot
EnvironmentFile=/home/rick/ozzy-bot/.env
ExecStart=/home/rick/ozzy-bot/venv/bin/python /home/rick/ozzy-bot/signal_generator.py
StandardOutput=journal
StandardError=journal
Environment="HERMES_WEBHOOK_URL=http://127.0.0.1:5001/webhook"
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PYTHONUNBUFFERED=1"
