# systemd/ozzybot-openclaw-market-sniper.service
[Unit]
Description=OzzyBot OpenClaw Market Sniper State Writer
After=network-online.target ozzybot-openclaw-macro-scout.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/rick/ozzy-bot
EnvironmentFile=/home/rick/ozzy-bot/config/live-micro.env
ExecStartPre=/usr/bin/mkdir -p /home/rick/ozzy-bot/shared
ExecStart=/home/rick/ozzy-bot/venv/bin/python /home/rick/ozzy-bot/core/market_sniper.py
StandardOutput=journal
StandardError=journal
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/rick/ozzy-bot"
